# praktikum-modul-1/task-3/crontabs
* * * * * /bin/bash /minute_log.sh

59 * * * * /bin/bash /aggregate_minutes_to_hourly_log.sh

59 23 * * * /bin/bash /backup_metrics.sh
