# praktikum-modul-1/task-2/crontabs
0 * * * * /home/oberon/Documents/ITS/MODUL_1/praktikum-modul-1-d27/task-2/yu_database.sh
